# services/hw_remote.service
[Unit]
Description=remote

[Service]
Type=simple
Restart=always
RestartSec=5
ExecStart=/usr/bin/python3 /home/spadia/mastpi-car/lumino_wifi_toolbox/hw_remote.py
ExecStop=/bin/systemctl kill hw_remote

[Install]
WantedBy=multi-user.target


 
