# services/masterpi.service
[Unit]
Description=Start MasterPi

[Service]
Type=simple
Restart=always
RestartSec=5
ExecStart=/home/spadia/lumino_bot/lumino/MasterPi.py 
StandardOutput=null
StandardError=null

[Install]
WantedBy=multi-user.target
